# services/nodejs/run
#!/bin/bash
export NODE_ENV=production
nodejs /var/www/index.js
